# packaging/debian.sa2dnsbl-daemon.postrm
#!/bin/sh

if [ "$1" = purge ]; then
   update-rc.d sa2dnsbld remove
fi
